# deploy/archlinux/rules.d/99-monitor-hotplug.rules
## reload bspwm when a screen is added or removed
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/usr/bin/bspc wm -r"
